# examples/xetra@.service
[Unit]
Description=Xetra %I Trade Data Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xetra
Group=xetra
WorkingDirectory=/var/lib/xetra
Environment="PATH=/opt/yf_parqed/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Run daemon mode with logging
# Note: Respects trading hours (08:30-18:00 CET/CEST) by default
# Add --active-hours "00:00-23:59" for 24/7 operation
ExecStart=/opt/yf_parqed/.venv/bin/xetra-parqed \
    --wrk-dir /var/lib/xetra \
    --log-file /var/log/xetra/%i.log \
    fetch-trades %I \
    --daemon \
    --interval 1 \
    --pid-file /run/xetra/%i.pid

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30

# Restart on failure
Restart=on-failure
RestartSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/xetra /var/log/xetra /run/xetra

# Create PID directory at startup
RuntimeDirectory=xetra
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
